# crash-mcp.service
[Unit]
Description=Crash MCP Server - Crash Dump Analysis Service
Documentation=https://github.com/42research/crash-mcp
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=crash-mcp
WorkingDirectory=/opt/crash-mcp

# Service execution
ExecStart=/usr/local/bin/crash-mcp-http
ExecReload=/bin/kill -HUP $MAINPID
ExecStop=/bin/kill -TERM $MAINPID

# Restart policy
Restart=on-failure
RestartSec=5s
StartLimitInterval=60s
StartLimitBurst=3

# Environment variables
Environment="PORT=8080"
Environment="LOG_LEVEL=INFO"
Environment="PYTHONUNBUFFERED=1"

# Tunnel integration - MCP_SERVER_URL will be set via systemctl set-environment or EnvironmentFile
# EnvironmentFile=/etc/crash-mcp/crash-mcp.env

# Process management
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=10s

# Security and isolation
# Note: Running as root is necessary because the 'crash' utility requires:
# - Access to kernel memory dumps in /var/crash
# - CAP_SYS_ADMIN capability to read kernel structures
# - Access to /dev/mem or /dev/kmem for live kernel analysis
# The service is still protected by systemd security features below
NoNewPrivileges=false
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/crash-dumps /var/log/crash-mcp /var/crash

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=crash-mcp

[Install]
WantedBy=multi-user.target

